# crawler/daily/crontab
* 9-10 * * * /usr/local/bin/python /corona-kakao-bot/daily.py > /proc/1/fd/1 2>/proc/1/fd/2
*/5 11-23 * * * /usr/local/bin/python /corona-kakao-bot/daily.py > /proc/1/fd/1 2>/proc/1/fd/2

